# machines/eb-jibri-template/etc/systemd/system/jibri.service.d/override.buster.conf
[Service]
Environment="PATH=/usr/lib/jvm/nvidia-java-8-openjdk-amd64/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/games:/usr/games"
