# machines/eb-jitsi/etc/systemd/system/nginx.service.d/override.conf
[Service]
TasksMax=65000
LimitNPROC=65000
LimitNOFILE=65000
